# conf/supervisor/zza.conf
[program:zza]
command=/usr/bin/node /home/pi/arduino_zza/web-ui/app.js
autorestart=true
user=root
directory=/home/pi/arduino_zza/web-ui
